# .config/autostart/xbindkeys.desktop
[Desktop Entry]
Exec=xbindkeys
Hidden=false
Icon=system-run
Path=
Terminal=false
Type=Application

