# .config/autostart/xembedsniproxy.desktop
[Desktop Entry]
Exec=xembedsniproxy
Icon=system-run
Path=
Terminal=false
Type=Application
