# systemd/user/lsyncd.service
[Unit]

[Service]
ExecStart=/usr/bin/lsyncd ${HOME}/.lsync/config.lua

[Install]
WantedBy=default.target
